# sound-greed.service
[Unit]
Description=sound greed service

[Service]
ExecStart=/path/to/bin/sound-greed

[Install]
WantedBy=default.target
